# config/supervisor.conf
[supervisord]
nodaemon=true

[program:server]
directory=/app
command=npm run start
autostart=true

[program:redis]
command=redis-server
autostart=true